# deploy/cron/mbprestige-scheduler.cron
* * * * * cd /var/www/mbprestige && php artisan schedule:run >> /dev/null 2>&1
